# rc.local
#!/bin/sh -e
sudo mount /dev/disk/by-id/nvme-Amazon_Elastic_Block_Store_vol0013c62fc213f8520 /mnt
sudo chown ubuntu.ubuntu /mnt
sudo swapon /swapfile
sudo swapon /swapfile2
sudo -u ubuntu nohup /home/ubuntu/continental-test/plot-generation/start-webserver.sh &
sudo -u ubuntu nohup /home/ubuntu/canadarasp/web-server/timezone/run-timezone-server.sh &
exit 0
